# sensor/hyperion-sensor.service
[Unit]
Description=Hyperion Solar Sensor Ingest
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/hyperion/sensor
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
LogRateLimitIntervalSec=30
LogRateLimitBurst=5

[Install]
WantedBy=multi-user.target
